# overlays/host/etc/systemd/system/stardock.service
[Unit]
Description=Stardock

[Service]
WorkingDirectory=/var/stardock
Environment="INITS=\
  /var/data/stardock-init-monitor-repos\
  /var/data/stardock-init-nginx\
  /var/data/stardock-init-nginx-hupper\
  /var/data/stardock-init-monitor-headquarters\
  /var/data/stardock-init-monitor-remotes\
  /var/data/stardock-init-reaper\
"
ExecStart=/bin/bash -c "/var/stardock/bin/stardock"
Restart=always

[Install]
WantedBy=multi-user.target
