# overlays/host/etc/cron.daily/restart_nginx
#!/bin/sh

/var/stardock/bin/stardock-restart-nginx
